# startup/systemd/devmail.service
[Unit]
Description=A small mail server for developers

[Service]
ExecStart=/usr/local/bin/devmail

[Install]
WantedBy=multi-user.target
